# files/etc/uci-defaults/disable-firewall.sh
#!/bin/sh
echo "Disabling firewall..."
/etc/init.d/firewall disable

